# scripts/etc/inittab
::sysinit:/etc/rcS
::respawn:-/bin/sh
tty2::askfirst:-/usr/bin/sh
tty3::askfirst:-/usr/bin/sh
::ctrlaltdel:/bin/backup

